# com.shadman.root-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  ROOT auto-refresh daemon for macOS.
  Replace ROOT_DIR with your actual root-kg install path, then:
    cp com.shadman.root-refresh.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.shadman.root-refresh.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.root-kg.refresh</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-c</string>
        <string>source ~/.zshrc 2>/dev/null; cd ROOT_DIR &amp;&amp; source .venv/bin/activate &amp;&amp; python indexer.py --extract 2>&amp;1</string>
    </array>

    <key>StartInterval</key>
    <integer>7200</integer>

    <key>StandardOutPath</key>
    <string>ROOT_DIR/logs/launchd.out.log</string>

    <key>StandardErrorPath</key>
    <string>ROOT_DIR/logs/launchd.err.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
